# headphone-event.service
[Unit]
Description=Headphone event handler

[Service]
ExecStart=/home/projedi/.bin/headphone-event

[Install]
WantedBy=default.target
